# com.pwagenet.claude-approval.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pwagenet.claude-approval</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/pwagenet/.bun/bin/bun</string>
    <string>/Users/pwagenet/Development/claude-approval-server/index.ts</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/claude-approval.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-approval.error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/pwagenet</string>
    <key>PATH</key>
    <string>/Users/pwagenet/.local/bin:/Users/pwagenet/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
